# cron.d/nightly
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
00 01  * * *   fewtarius   ~/repositories/fewtarius/351ELEC-builder/builder >~/repositories/fewtarius/351ELEC-builder/builder.log
